# container-update-db-lop/crontab
* * * * * python /update_lop_db.py >> /logfile

